# edge/deploy/cryptiq-wg-sync.service
[Unit]
Description=Sync CryptiQ WireGuard peers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cryptiq-wg-sync
